# mdkupdate.cron
#!/bin/bash

urpmi.update -a
urpmi --force --auto --auto-select


