# rdate.cron
#!/bin/bash

rdate -s time.nrc.ca
hwclock --utc --systohc

